# etc/cron.hourly/furneaux-sync
#!/bin/sh
exec furneaux >> /var/log/furneaux.log
